# supervisor/app.conf
[program:app]
priority=5
directory=/app/code
environment=PORT=3000,HOSTNAME=0.0.0.0
command=node -r dotenv-flow/config --expose-gc /app/code/app/dist/server/app.js
user=cloudron
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
